# xivo-call-logs/debian/cron.d
#
# cron jobs for xivo-call-logs
#

*/5 * * * * root /usr/bin/xivo-call-logs
