# xivo-call-logs/debian/cron.d
#
# cron jobs for xivo-call-logs
#

25 4 * * * root /usr/bin/xivo-call-logs
