# debian/cron.d
#
# cron jobs for wazo-call-logs
#

25 4 * * * root /usr/bin/wazo-call-logs
